# etc/supervisor/conf.d/co2meter.ini
[program:co2meter]
directory = /tmp
command = /usr/bin/python /home/pi/co2meter/main.py
stdout_logfile = /var/log/supervisor/co2meter.log
stderr_logfile = /var/log/supervisor/co2meter.log
autostart=true
autorestart=true
startsecs=10
